# feelworld.service
[Unit]
Description=FEELWORLD-05 teleprompter controller
After=bluetooth.target
Wants=bluetooth.target

[Service]
Type=simple
WorkingDirectory=%h/pilot_feelworld_as_keyboard
ExecStart=%h/pilot_feelworld_as_keyboard/venv/bin/python controller.py
Environment=PYTHONUNBUFFERED=1
StandardOutput=journal
StandardError=journal
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
